# build/ansible/roles/netfilter/files/netfilter.service
[Unit]
Description=IPv4 NetFilter with iptables & ipset
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/netfilter.sh reload
ExecStartPost=/bin/echo "NetFilter [STARTED]"

TimeoutStartSec=60

[Install]
WantedBy=basic.target
